# systemd/system/etc-iscsi-send_targets.mount
[Unit]
Description=iSCSI send_targets record-keeping

[Mount]
What=tmpfs
Where=/etc/iscsi/send_targets
Type=tmpfs
Options=mode=0700,nodev