# debian/cluster-netboot.etc-iscsi-nodes.mount
[Unit]
Description=iSCSI nodes record-keeping

[Mount]
What=tmpfs
Where=/etc/iscsi/nodes
Type=tmpfs
Options=mode=0700,nodev