# systemd/system/systemd-networkd-wait-online.service.d/minimal-operational.conf
# Explicitly wait for eth0 to be routable, but only for 15 seconds
[Service]
ExecStart=
ExecStart=/lib/systemd/systemd-networkd-wait-online --timeout=15 --interface=eth0:routable